# config/rpi/conf/homer1.service
[Unit]
Description=Homer1

[Service]
ExecStart=/usr/bin/docker-compose -f /root/docker-compose.yml up

[Install]
WantedBy=multi-user.target

